# config/systemd/user/backup.service
[Unit]
Description=Run obnam backup

[Service]
Type=simple
ExecStart=/usr/bin/obnam backup
ExecStartPost=/usr/bin/obnam forget
